# flask-elastic-search/develop/supervisord.conf
[supervisord]

[program:elasticsearch]
command = elasticsearch-6.3.2/bin/elasticsearch
stopasgroup = true
autorestart = true
